# production-audit/server-snapshot/systemd/gpt56-retry-proxy.service.d/concurrency.conf
[Service]
ExecStart=
ExecStart=/usr/bin/python3 /opt/gpt56-retry-proxy/gpt56_retry_proxy.py --listen 0.0.0.0:3000 --upstream http://127.0.0.1:3001 --attempts 5 --upstream-concurrency 800
